# etc/flux-coral2-dws.service.in
[Unit]
Description=Flux-DWS communication service
BindsTo=flux.service
After=flux.service

[Service]
TimeoutStopSec=90
KillMode=mixed
ExecStart=@X_BINDIR@/flux python @X_BINDIR@/coral2_dws.py
SyslogIdentifier=flux-coral2-dws
Restart=always
RestartSec=10s
RestartPreventExitStatus=2 3
WatchDogSec=60s

User=flux
Group=flux

[Install]
WantedBy=flux.service
